# backend/proton.mount
[Mount]
Type=proton
Exec=@libexecdir@/gvfsd-proton
AutoMount=false
Scheme=proton
